# contrib/systemd/fiche.socket
[Unit]
Description=fiche

[Socket]
BindIPv6Only=both
ListenStream=99

[Install]
WantedBy=sockets.target
